# feeds/luci/applications/luci-quick/root/etc/crontabs/root
*/3 * * * * sh /usr/sbin/g3monitor
*/10 * * * * sh /usr/sbin/savelog
0 * * * * /bin/EV_server
