# package/firmware-ev/files/etc/init.d/ev_boot
#!/bin/sh /etc/rc.common
START=13

start() {
	/bin/ev init_boot
	/bin/dashboard mqtt_sub
}
